# .config/autostart/emacs-autostart.desktop
[Desktop Entry]
Type=Application
Name=Emacs Autostart
Exec=env emacs -f my/configure-exwm &
Terminal=false
