# .config/autostart/exwm-autostart.desktop
[Desktop Entry]
Type=Application
Name=Emacs Autostart
Exec=env emacs -f pjs-configure-exwm
Terminal=false
