# .config/autostart/exwm-autostart.desktop
[Desktop Entry]
Type=Application
Name=Emacs Autostart
Exec=/home/paul/bin/exwm
Terminal=false
