# .config/autostart/exwm-autostart.desktop
[Desktop Entry]
Type=Application
Name=Emacs Autostart
Exec=exwm
Terminal=false
